# fs/shared/usr/share/initramfs-tools/crontab
* * * * * sleep 20 ; /usr/local/sbin/telegram-notify.sh >> /tmp/telegram-notify.cron.log 2>&1 &